# straper/db/public/services/custom-units/docker.service.d/override.conf
[Unit]
After=nftables.service
Wants=nftables.service
